# packaging/debian/tagstackd.service
[Unit]
Description=Mcaster1 TagStack daemon
After=network-online.target mariadb.service
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/tagstack.mcaster1.com/current
ExecStart=/var/www/tagstack.mcaster1.com/current/daemon/build/tagstackd --config /var/www/tagstack.mcaster1.com/shared/etc/tagstackd.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=3
StandardOutput=append:/var/www/tagstack.mcaster1.com/shared/log/tagstackd.stdout
StandardError=append:/var/www/tagstack.mcaster1.com/shared/log/tagstackd.stderr

# Capability so www-data can bind 9890 (well above 1024 so not strictly needed,
# but explicit for parity with daemons that bind <1024)
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Sandboxing
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/www/tagstack.mcaster1.com/shared/log /var/www/tagstack.mcaster1.com/shared/cache
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
